# systemd/.config/systemd/user/wallpaper-switch.timer
[Unit]
Description=Run wallpaper switch periodically
After=graphical.target

[Timer]
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
